# target/linux/ky/riscv64/base-files/etc/hotplug.d/net/40-net-smp-affinity
#!/bin/sh

[ "$ACTION" = add ] || exit

get_device_irq() {
	local device="$1"
	local line
	local seconds="0"

	# wait up to 10 seconds for the irq/device to appear
	while [ "${seconds}" -le 10 ]; do
		line=$(grep -E -m 1 "${device}\$" /proc/interrupts) && break
		seconds="$(( seconds + 2 ))"
		sleep 2
	done
	echo ${line} | sed 's/:.*//'
}

set_interface_core() {
	local core_mask="$1"
	local interface="$2"
	local device="$3"

	[ -z "${device}" ] && device="$interface"

	local irq=$(get_device_irq "$device")

	echo "${core_mask}" > /proc/irq/${irq}/smp_affinity
}

case "$(board_name)" in
armsom,sige7|\
friendlyarm,nanopc-t6|\
friendlyarm,nanopi-r3s|\
friendlyarm,nanopi-r5c|\
friendlyarm,nanopi-r6c|\
radxa,e25|\
sinovoip,rk3568-bpi-r2pro)
	set_interface_core 2 "eth0"
	set_interface_core 4 "eth1"
	;;
friendlyarm,nanopi-r2c|\
friendlyarm,nanopi-r2c-plus|\
friendlyarm,nanopi-r2s|\
radxa,cm3-io|\
xunlong,orangepi-r1-plus|\
xunlong,orangepi-r1-plus-lts)
	set_interface_core 2 "eth0"
	set_interface_core 4 "eth1" "xhci-hcd:usb[0-9]+"
	;;
friendlyarm,nanopi-r4s|\
friendlyarm,nanopi-r4s-enterprise)
	set_interface_core 10 "eth0"
	set_interface_core 20 "eth1"
	;;
friendlyarm,nanopi-r5s|\
friendlyarm,nanopi-r6s)
	set_interface_core 2 "eth0"
	set_interface_core 4 "eth1"
	set_interface_core 8 "eth2"
	;;
ky,orangepi-rv2)
	echo 1 > /sys/class/wlan_power_ctrl/wlan_power/wlan_power
	;;
ky,orangepi-r2s)
	# 8-core CPU: use all cores for network processing
	# eth2 = LAN (2.5G RTL8125), eth3 = WAN (2.5G RTL8125)

	# Set IRQ affinity for RTL8125 queues
	# eth2 (LAN) -> CPU 0-3 (mask: 0f)
	# eth3 (WAN) -> CPU 4-7 (mask: f0)
	for irq in $(grep -E "eth2-" /proc/interrupts | awk '{print $1}' | tr -d ':'); do
		[ -n "$irq" ] && echo 0f > /proc/irq/$irq/smp_affinity 2>/dev/null
	done
	for irq in $(grep -E "eth3-" /proc/interrupts | awk '{print $1}' | tr -d ':'); do
		[ -n "$irq" ] && echo f0 > /proc/irq/$irq/smp_affinity 2>/dev/null
	done

	# Enable RPS (Receive Packet Steering) on all RX queues
	# Distribute to all 8 cores (mask: ff)
	for q in /sys/class/net/eth2/queues/rx-*; do
		[ -d "$q" ] && echo ff > "$q/rps_cpus" 2>/dev/null
	done
	for q in /sys/class/net/eth3/queues/rx-*; do
		[ -d "$q" ] && echo ff > "$q/rps_cpus" 2>/dev/null
	done

	# Set RPS flow count for better flow distribution
	for q in /sys/class/net/eth2/queues/rx-*; do
		[ -d "$q" ] && echo 4096 > "$q/rps_flow_cnt" 2>/dev/null
	done
	for q in /sys/class/net/eth3/queues/rx-*; do
		[ -d "$q" ] && echo 4096 > "$q/rps_flow_cnt" 2>/dev/null
	done

	# Enable XPS (Transmit Packet Steering) on all TX queues
	for q in /sys/class/net/eth2/queues/tx-*; do
		[ -d "$q" ] && echo ff > "$q/xps_cpus" 2>/dev/null
	done
	for q in /sys/class/net/eth3/queues/tx-*; do
		[ -d "$q" ] && echo ff > "$q/xps_cpus" 2>/dev/null
	done

	# Set global RPS flow entries
	echo 32768 > /proc/sys/net/core/rps_sock_flow_entries 2>/dev/null
	;;
esac

